# ops/systemd/myhibachi-frontend.service
[Unit]
Description=My Hibachi Frontend Service
After=network.target
Requires=network.target

[Service]
Type=exec
User=myhibachi
Group=myhibachi
WorkingDirectory=/opt/myhibachi/myhibachi-frontend
Environment=PATH=/opt/myhibachi/node/bin:/usr/local/bin:/usr/bin:/bin
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/opt/myhibachi/node/bin/npm start
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myhibachi-frontend

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/myhibachi/logs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=2048
MemoryMax=1G
CPUQuota=150%

[Install]
WantedBy=multi-user.target