# ops/systemd/myhibachi-backup.service
[Unit]
Description=My Hibachi Database Backup Service
After=postgresql.service
Wants=postgresql.service

[Service]
Type=oneshot
User=myhibachi
Group=myhibachi
Environment=PATH=/opt/myhibachi/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=ENV_FILE=/opt/myhibachi/config/.env
ExecStart=/opt/myhibachi/venv/bin/python /opt/myhibachi/ops/backup_db.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=myhibachi-backup

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/backups/myhibachi /var/log/myhibachi
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true