# systemd/myhibachi-api.service
[Unit]
Description=MyHibachi FastAPI (main)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/srv/myhibachi/api
Environment="PORT=8001"
Environment="ENV=production"
ExecStart=/srv/myhibachi/api/.venv/bin/gunicorn app.main:app -k uvicorn.workers.UvicornWorker -b 127.0.0.1:8001 --workers 2 --threads 2 --timeout 60
Restart=always

[Install]
WantedBy=multi-user.target
