# infrastructure/systemd/myhibachi-api.service
[Unit]
Description=My Hibachi Chef Unified API
Documentation=https://github.com/suryadizhang/mh-project-
After=network.target postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/vhosts/api.myhibachichef.com/myhibachi-api

# Environment configuration
Environment="PATH=/var/www/vhosts/api.myhibachichef.com/myhibachi-api/venv/bin"
EnvironmentFile=/var/www/vhosts/api.myhibachichef.com/myhibachi-api/.env

# Gunicorn configuration for unified API
ExecStart=/var/www/vhosts/api.myhibachichef.com/myhibachi-api/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 300 \
    --keepalive 2 \
    --access-logfile logs/access.log \
    --error-logfile logs/error.log \
    --log-level info \
    --capture-output \
    --enable-stdio-inheritance \
    src.main:app

# Pre/post execution scripts
ExecStartPre=/bin/mkdir -p /var/www/vhosts/api.myhibachichef.com/myhibachi-api/logs
ExecReload=/bin/kill -s HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Resource limits (optimized for 4GB VPS)
MemoryLimit=3G
MemoryAccounting=yes
CPUQuota=150%
TasksMax=4096

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/www/vhosts/api.myhibachichef.com/myhibachi-api/logs
ReadWritePaths=/var/www/vhosts/api.myhibachichef.com/myhibachi-api/uploads

# Additional security
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target